# files/etc/inittab
::sysinit:/etc/init.d/rcS S boot
::shutdown:/etc/init.d/rcS K shutdown
ttyS0::askfirst:/usr/libexec/login.sh
tty1::respawn:/usr/bin/miner_start.sh
tty2::askfirst:/usr/libexec/login.sh
